# config/uwsgi.silica.conf
description "uwsgi processes"
start on runlevel [2345]
stop on runlevel [06]

exec /home/xangis/.virtualenvs/silica/bin/uwsgi --ini /var/django/silica/uwsgi.ini
